# .config/systemd/user/ghostmirror.path
[Path]
PathChanged=%h/.config/ghostmirror/mirrorlist.all

