# .config/systemd/user/mbsync.path
[Path]
PathChanged=%h/Mail/.notmuch

[Install]
WantedBy=default.target
