# .config/systemd/user/timecheck.service
[Unit]
Description=Sends heads-up notifications every 15 minutes.

[Service]
Type=oneshot
ExecStart=python %h/.local/bin/timecheck.py
