# .config/systemd/user/timecheck.timer
[Unit]
Description=Sends heads-up notifications every 15 minutes.

[Timer]
OnActiveSec=1
OnUnitActiveSec=15min
AccuracySec=1us
