# .config/systemd/user/trash-empty.service
[Unit]
Description=Empty trash
Documentation=https://github.com/andreafrancia/trash-cli

[Service]
Type=oneshot
ExecStart=trash-empty -f 30
