# .config/systemd/user/elephant.service.d/overrides.conf
[Unit]
ConditionEnvironment=WAYLAND_DISPLAY

[Service]
Slice=session-graphical.slice
