# .config/systemd/user/hypridle.service.d/overrides.conf
[Service]
ExecCondition=/lib/systemd/systemd-xdg-autostart-condition "wlroots:sway:Wayfire:labwc:Hyprland" ""
Slice=background-graphical.slice
