# .config/systemd/user/vdirsyncer.timer.d/overrides.conf
[Timer]
OnBootSec=
OnUnitActiveSec=
OnBootSec=5m
OnUnitActiveSec=60m

