# .config/systemd/user/hyprpolkitagent.service.d/override.conf
[Unit]
Before=waybar.service

[Service]
Type=exec
ExecStartPost=sleep 0.5
